# hamr.service
[Unit]
Description=Hamr Launcher for Quickshell
Documentation=https://github.com/Stewart86/hamr
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/qs -c hamr
Restart=on-failure
RestartSec=3

[Install]
WantedBy=niri.service
